# server/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tomcat]
command=node /app/appserver/appserver.js
autorestart=true
